# antigravity-terminal/rootfs/etc/s6-overlay/s6-rc.d/antigravity-terminal/run
#!/command/with-contenv bashio
# s6-overlay service: antigravity-terminal
# This is the main process that s6 monitors and restarts if needed

exec /usr/bin/antigravity-startup
